# debian/statsd-c.statsd.upstart.in
# statsd-c - line-compatible statsd server in c
#

description	"statsd daemon"

start on filesystem or runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5
umask 022

pre-start script
    mkdir -p -m0770 /var/run/statsd
    chown statsd:statsd /var/run/statsd

    mkdir -p -m0770 /var/lib/statsd
    chown statsd:statsd /var/lib/statsd
end script

script
	DAEMON_ARGS
	[ -e /etc/default/statsd ] && . /etc/default/statsd
	exec su -s /bin/sh -c "/usr/sbin/statsd ${DAEMON_ARGS} -l /var/run/statsd -D" statsd
end script
