# debian/statsd-c.statsd.upstart
# statsd-c - line-compatible statsd server in c
#

description	"statsd daemon"

start on filesystem or runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5
umask 022

pre-start script
    test -x /var/lib/statsd || { mkdir -p /var/lib/statsd; chown statsd:statsd /var/lib/statsd-c }
    mkdir -p -m0755 /var/run/statsd
end script

script
	[ -e /etc/default/statsd ] && . /etc/default/statsd
	exec /usr/sbin/statsd ${DAEMON_ARGS} -l /var/run/statsd-c -D
end script
